# misc/configs/supervisor.conf
[program:apiserver]
directory=/path/to/project
command=/path/to/project/app
process_name=%(program_name)s
redirect_stderr=false
stdout_logfile=/path/to/project/log/%(program_name)s.log
stderr_logfile=/path/to/project/log/%(program_name)s.log
autorestart=true
stdout_logfile_maxbytes=500MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
stdout_capture_maxbytes=500MB   ; number of bytes in 'capturemode' (default 0)
stdout_events_enabled=true   ; emit events on stdout writes (default false)
stderr_logfile_maxbytes=500MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)
stderr_capture_maxbytes=500MB   ; number of bytes in 'capturemode' (default 0)
stderr_events_enabled=true   ; emit events on stderr writes (default false)
